# deploy/systemd/fossgis-platzmonitor-backend.service
[Unit]
Description=FOSSGIS Platzmonitor Backend (FastAPI / uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Optional: fuer Systemdienst setzen. Fuer User Dienst muessen diese Zeilen entfernt bleiben.
# User=www-data
# Group=www-data

# Anpassen: absoluter Pfad zum backend/-Verzeichnis des geklonten Repos
WorkingDirectory=/path/to/platzmonitor/backend

Environment=PYTHONPATH=src
# Laedt backend/.env (optional; Zeile mit - = Datei darf fehlen)
EnvironmentFile=-/path/to/platzmonitor/backend/.env

# Wie deploy/nginx location /event-api/: Swagger UI muss /event-api/openapi.json laden, nicht /openapi.json.
# Nur FastAPI-ROOT_PATH (siehe app), nicht uvicorn --root-path (sonst stimmt path-Routing mit nginx-rewrite nicht).
Environment=ROOT_PATH=/event-api

# uv muss im PATH liegen (z. B. Symlink nach /usr/local/bin oder PATH unten erweitern)
# Environment=PATH=/usr/local/bin:/usr/bin

# Kein Access-Log hier: wiederholte GETs spammen das Journal; Zugriffe siehe nginx access.log
ExecStart=/usr/bin/env uv run uvicorn app.main:app --host 127.0.0.1 --port 8000 --no-access-log
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
